# conf/crontab
*/1 * * * * /Library/PostgreSQL/9.1/bin/psql -U pglog -c "SELECT PgLog_Import()" pglog
